# webapp/crontab
*/1 * * * * rails crawl > /var/log/cron.log
